# shadowsocks.conf
[program:shadowsocks] 
command=python /root/shadowsocks/server.py 
autorestart=true 
user=root 
